# docker/supervisord.conf
[program:nicistore]
command=node index.js
directory=/app
autostart=true
autorestart=true
environment=NODE_ENV=production,DEBUG=nicistore
stderr_logfile=/var/log/nicistore.err.log
stdout_logfile=/var/log/nicistore.out.log
